# systemd/ecoflow-nut-bridge.service
[Unit]
Description=EcoFlow NUT Bridge
After=network.target bluetooth.target
Wants=bluetooth.target
# Start before NUT so the dummy-ups state file is seeded (ExecStartPre below)
# before the driver reads it. install.sh also adds a matching After= drop-in to
# nut-server so the ordering holds from both sides.
Before=nut-server.service

[Service]
Type=simple
User=ecoflow
Group=nut
WorkingDirectory=/opt/ecoflow-nut-bridge
# Make Bluetooth reboot-safe: unblock rfkill and power on the adapter before we
# scan. Some boots come up soft-blocked / unpowered. The leading "+" runs this
# as root (the service itself runs as the unprivileged ecoflow user).
ExecStartPre=+/bin/sh -c 'rfkill unblock bluetooth; sleep 1; bluetoothctl power on >/dev/null 2>&1 || true'
# Seed an initial state file so the dummy-ups driver has something to read at
# boot. Best-effort ("-"): a hiccup here must not stop the bridge itself.
ExecStartPre=-/opt/ecoflow-nut-bridge/.venv/bin/ecoflow-nut --config /etc/ecoflow-nut/config.yaml seed
ExecStart=/opt/ecoflow-nut-bridge/.venv/bin/ecoflow-nut --config /etc/ecoflow-nut/config.yaml run
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
# Persistent state dir for the optional local SQLite telemetry store. systemd
# creates /var/lib/ecoflow-nut owned by the service user (the default
# sqlite.path lives here). Harmless if SQLite logging is left disabled.
StateDirectory=ecoflow-nut

# Note: /run/nut is created and owned by the nut-server package (nut:nut, 0770);
# the ecoflow user is added to the nut group so it can write the state file.
# We deliberately do NOT use RuntimeDirectory=nut here -- that would delete
# /run/nut when the bridge restarts and break a running upsd.

[Install]
WantedBy=multi-user.target
